# youtube-agent/com.marcbyers.youtube-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marcbyers.youtube-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/marcbyers/Agents/youtube-agent/venv/bin/python</string>
        <string>/Users/marcbyers/Agents/youtube-agent/agent.py</string>
        <string>--daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/marcbyers/Agents/youtube-agent</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/marcbyers/Agents/youtube-agent/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/marcbyers/Agents/youtube-agent/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/marcbyers/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
